# hosts/common/optional/podman/homeassistant-pip/etc/s6-overlay/s6-rc.d/setBluetoothOneshot/up
foreground { echo 'Setting up bluetooth' }
foreground { /usr/local/bin/setBluetoothCapabilities.sh }
